# root/app/qbittorrent/supervisord.conf
[program:qbittorrent]
command = /usr/bin/qbittorrent-nox
autorestart = true
priority = 1
startsecs = 3
stdout_logfile=/var/log/qbittorrent.log
logfile_maxbytes = 50MB
